# crontab/root
0 4 * * 1	/usr/local/sbin/update_script.sh >>/var/log/update_script.log 2>&1
@reboot		/usr/local/sbin/update_script.sh >>/var/log/update_script.log 2>&1
0 0 * * *	/usr/local/sbin/watch_cron.sh
