# crontab/update_script.sh
#!/bin/sh

date
apt --yes update && apt --yes upgrade
