# crontab/update_script.sh
#!/bin/sh

apt --yes update && apt --yes upgrade
